# HOTEL/KICKSTART/autostart/disable_shortcut_2.desktop
[Desktop Entry]
Exec=xfconf-query -c xfce4-keyboard-shortcuts -p '/commands/custom/<Alt>F3' -r -R
Terminal=false