# HOTEL/KICKSTART/autostart/office_startup.desktop
[Desktop Entry]
Exec=bash -c "xfce4-terminal -H -e 'sh /usr/local/bin/office_instructions.sh'"
Terminal=false