# HOTEL/KICKSTART/autostart/office_startup.desktop
[Desktop Entry]
Exec=bash -c "sleep 5 && xfce4-terminal -H -e 'sudo /usr/local/bin/office_instructions.sh'"
Terminal=false