# KICKSTART/autostart/autores.desktop
[Desktop Entry]
Exec=/usr/local/bin/autores.sh
Terminal=false