# KICKSTART/autostart/check.desktop
[Desktop Entry]
Exec=bash -c "sleep 5 && xfce4-terminal -H -e 'cat /tmp/script_log.log'"
Terminal=false