# KICKSTART/autostart/delete.desktop
[Desktop Entry]
Exec=bash -c "sleep 10 && /usr/local/bin/check.sh"
Terminal=false