# KICKSTART/autostart/disable_desktop_icons.desktop
[Desktop Entry]
Exec=xfconf-query -c xfce4-desktop -np /desktop-icons/style -t 'int' -s '0'
