# KICKSTART/autostart/pacmd_volume.desktop
[Desktop Entry]
Exec=bash -c "sleep 5 && /usr/local/bin/pacmd_volume.sh"
Terminal=false