# KICKSTART/autostart/pacmd.desktop
[Desktop Entry]
Exec=bash -c "sleep 5 && /usr/local/bin/pacmd.sh"
Terminal=false