# KICKSTART/autostart/pacmd_disablepci.desktop
[Desktop Entry]
Exec=bash -c "sleep 5 && /usr/local/bin/pacmd_disablepci.sh"
Terminal=false